# scripts/maintenance/maintenance.crontab
# System maintenance cron job
# Runs every Tuesday at 11:00 AM
0 11 * * 2 $HOME/.config/scripts/maintenance/maintenance.sh >> /tmp/maintenance.log 2>&1
